# services/tomcat/run
#!/usr/bin/env bash

exec /sbin/setuser ${TOMCAT_USER} ${TOMCAT_BIN_PATH}/catalina.sh run
